# infra/systemd/blockbook-bcash.service
[Unit]
Description=BlockBook BCH indexer (cashtokens-aware)
After=bchn.service network.target
Wants=bchn.service
Documentation=https://github.com/mainnet-pat/blockbook

[Service]
Type=simple
User=blockbook-bcash
Group=blockbook-bcash
WorkingDirectory=/opt/blockbook-src
ExecStart=/opt/blockbook-src/blockbook \
    -blockchaincfg=/etc/blockbook-bcash/blockchaincfg.json \
    -datadir=/var/lib/blockbook-bcash/db \
    -sync \
    -internal=127.0.0.1:9031 \
    -public=127.0.0.1:9131 \
    -certfile= \
    -logtostderr
Restart=on-failure
RestartSec=30s

# TimeoutStopSec=900 gives RocksDB room to flush its write-ahead log and close
# column families cleanly during shutdown. Previous 600s killed mid-flush
# once under load; raising it costs nothing when shutdown is fast.
TimeoutStopSec=900s

# Memory guardrails. 16 GB box; BCHN + Postgres + SvelteKit + workers need
# ~5-6 GB combined in steady-state. BlockBook's IBD-phase working set grows
# to ~10 GB if uncapped (observed: OOM-killed at anon-rss 10,697,728 kB
# during a deploy that also ran `cargo build --release` concurrently).
#
# First attempt at MemoryHigh=6G / MemoryMax=7G throttled RocksDB's block
# cache hard enough that IBD dropped from ~650 blocks/min to ~72 blocks/min
# (6-day ETA vs ~24h). The bottleneck was cache thrashing under cgroup
# reclaim pressure, not CPU. So: raise caps to accommodate the observed
# 10G working set, and rely on the install doc's "don't deploy during IBD"
# advisory + OOMPolicy=stop for protection instead of starving the cache.
#
# MemoryHigh is a soft throttle — cgroup v2 applies reclaim pressure once
# exceeded, but doesn't kill. 8 GB = BlockBook's happy working set without
# squeezing the block cache.
#
# MemoryMax is the hard kill line. 10 GB caps the observed pre-OOM peak
# and still leaves ~6 GB for everything else on the box.
#
# OOMPolicy=stop means systemd's own OOM handler will stop() the unit
# rather than letting RestartSec loop us into rapid crash cycles if we
# do hit MemoryMax.
MemoryHigh=8G
MemoryMax=10G
OOMPolicy=stop

# Hardening profile matching sync-tail / sync-cauldron / tokenstork.service.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/blockbook-bcash
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
