# infra/systemd/briefing.service
[Unit]
Description=Stork Sightings daily BCH token briefing generator
Documentation=https://github.com/Panmoni/tokenstork
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=tokenstork
Group=tokenstork
WorkingDirectory=/opt/tokenstork
# Absolute path to tsx in node_modules, matching the sync-leaderboards
# pattern. No PATH dependency on the tokenstork user's shell profile.
ExecStart=/opt/tokenstork/node_modules/.bin/tsx /opt/tokenstork/scripts/briefing.ts
EnvironmentFile=-/etc/tokenstork/env
Environment="BRIEFING_OUTPUT_DIR=briefings"
Environment="BRIEFING_PUBLIC_URL=https://tokenstork.com"

# ReadWritePaths lets the service write briefing files to the persistent
# briefings/ directory (at /opt/tokenstork/briefings). The SvelteKit app
# serves them from the same path.
ReadWritePaths=/opt/tokenstork/briefings

# Hardening — same profile as sync-leaderboards / sync-cauldron-stats.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
