# infra/systemd/enrich-seed.service
[Unit]
Description=tokenstork enrich-seed — one-shot bootstrap of live_token_utxo from BlockBook
After=network.target postgresql.service blockbook.service sync-enrich.service sync-bcmr-onchain.service
Conflicts=sync-enrich.service sync-bcmr-onchain.service
Wants=postgresql.service blockbook.service

[Service]
Type=oneshot
User=tokenstork
Group=tokenstork
WorkingDirectory=/opt/tokenstork
ExecStart=/opt/tokenstork/workers/target/release/enrich-seed

Environment="RUST_LOG=info"
# `-` prefix: start even if /etc/tokenstork/env is missing. The binary will
# then exit with an anyhow chain naming the first missing env var — operator
# sees the failure in journalctl.
EnvironmentFile=-/etc/tokenstork/env

# CPU containment — a full 17.5k-category sweep pins BlockBook hard.
# 200% (2 cores) keeps the box responsive during the seed run.
CPUQuota=200%

# Hardening — same profile as sync-enrich / sync-bcmr-onchain / sync-verify / sync-tail.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallArchitectures=native
UMask=0077
