# infra/systemd/sync-bcmr.timer
[Unit]
Description=Run tokenstork BCMR hydrator hourly

[Timer]
# First run 5 min after boot (give Postgres time to settle), then
# hourly. The worker is stale-first so steady-state load against
# Paytaca is capped by CAULDRON_MAX_RPS in the env; new categories +
# freshly-updated manifests hydrate within an hour of landing.
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
RandomizedDelaySec=1min
Unit=sync-bcmr.service

[Install]
WantedBy=timers.target
