# infra/systemd/sync-cauldron.timer
[Unit]
Description=Run tokenstork Cauldron full-scan every 4 hours (discovery + prune)

[Timer]
# Full-scan pass: walk every FT / FT+NFT category to (a) discover
# newly-listed tokens that weren't in the last fast-pass set and
# (b) prune tokens that have dropped off Cauldron.
#
# Runtime ~15 min at CAULDRON_MAX_RPS=5 due to the ~3,350 price
# checks. Four hours between runs gives disk + Postgres breathing
# room alongside BlockBook IBD + sync-tail + the fast-pass timer.
# Fresh prices for already-listed tokens are supplied by the
# sync-cauldron-fast.timer companion (every 10 min, ~60-90s/run).
OnBootSec=15min
OnUnitActiveSec=4h
Persistent=true
RandomizedDelaySec=5min
Unit=sync-cauldron.service

[Install]
WantedBy=timers.target
